# files/gnomeos/generate-initramfs/modules/50-cdrom/sr-loop@.service
[Unit]
DefaultDependencies=no
ConditionKernelCommandLine=root=live:gnomeos

[Service]
Type=oneshot
ExecStart=/usr/bin/losetup -f -P /dev/%I
