# files/gnomeos/generate-initramfs/modules/50-cdrom/sr-loop-partscan@.service
[Unit]
DefaultDependencies=no
ConditionKernelCommandLine=root=live:gnomeos

[Service]
Type=oneshot
ExecStart=/usr/bin/partx -u /dev/%I
