# files/gnomeos/generate-initramfs/modules/50-subpartitions/subpartitions@.service
[Unit]
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/bin/losetup -P -f /dev/%I
