# files/gnomeos/iso-extensions/run-mount-eoslive.mount
[Mount]
What=LABEL=EOSLIVE
Where=/run/mount/eoslive
Type=vfat
Options=ro,noatime,nodev,nosuid,noexec,nofail
DirectoryMode=0755
